# deploy/personal-blog.service
[Unit]
Description=personal-blog Next.js app
After=network.target

[Service]
Type=simple
User=cjl
WorkingDirectory=/home/cjl/projects/personal-blog
Environment=NODE_ENV=production
Environment=PORT=8080
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
